# recipes-core/server/tapecontrol-0.1/tapecontrol.service
[Unit]
Description=Raspebrry Pi Tape Remote TapeControl
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tapecontrol

[Install]
WantedBy=multi-user.target
